# systemd/nac-tap.service
[Unit]
Description=NAC-Tap - Transparent Network Bridge Monitor
Documentation=https://github.com/your-repo/nac-tap
After=network-online.target wlan0-up.service hostapd.service dnsmasq.service
Wants=network-online.target
Requires=hostapd.service wlan0-up.service

[Service]
Type=simple
User=root
# WorkingDirectory and ExecStart will be auto-updated by install-service.sh
WorkingDirectory=/path/to/nac-tap
ExecStart=/usr/bin/python3 /path/to/nac-tap/nac-tap.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nac-tap

# Security settings
NoNewPrivileges=false
PrivateTmp=false

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

